# dot_config/systemd/user/Gdrive.timer
[Unit] 
Description=GoogleDrive periodic sync
[Timer]
OnUnitInactiveSec=1 min
Persistent=True
[Install]
WantedBy=timers.target
