# dot_config/systemd/user/Gdrive.path
[Unit] 
Description=Google Drive folder watch
[Path]
PathExistsGlob=%h/Sync/GoogleDrive/*
[Install]
WantedBy=paths.target
