# dot_config/systemd/user/OneDrive.timer
[Unit] 
Description=OneDrive periodic sync
[Timer]
OnUnitInactiveSec=1 min
Persistent=True
[Install]
WantedBy=timers.target
