# dot_config/systemd/user/OneDrive.path
[Unit] 
Description=OneDrive folder watch
[Path]
PathExistsGlob=%h/Sync/OneDrive/*
[Install]
WantedBy=paths.target
