# dot_config/profile.d/10-display.sh.tmpl
{{ if (contains "WSL" .chezmoi.kernel.osrelease) }}#!/usr/bin/env sh

export DISPLAY=$(cat /etc/resolv.conf | grep nameserver | awk '{print $2}'):0
{{ end }}
